# rootfs/etc/services.d/shuffle-bot/finish
#!/usr/bin/execlineb -S1

# only tell s6 to bring down the entire container, if it isn't already doing so
# http://skarnet.org/software/s6/s6-supervise.html
if { s6-test ${1} -ne 0 }
if { s6-test ${1} -ne 256 }

# https://github.com/just-containers/s6-overlay/issues/193#issuecomment-263316923
redirfd -w 2 /dev/null s6-svscanctl -t /var/run/s6/services
